# framesync/framesync.timer
[Unit]
Description=Run OneDrive photo sync nightly

[Timer]
OnCalendar=hourly
AccuracySec=1m
Persistent=true

[Install]
WantedBy=timers.target
